# colreg-ui/src/deb/control/postinst
#!/bin/bash
# invoke-rc.d tomcat7 start
# init mongo db
# FIXME do we need any checks here before adding this data to the database?
mongoimport /var/dfa/colreg-ui/data/accessType.json
mongoimport /var/dfa/colreg-ui/data/accrualPeriodicity.json
mongoimport /var/dfa/colreg-ui/data/agentRelationType.json
mongoimport /var/dfa/colreg-ui/data/itemType.json
mongoimport /var/dfa/colreg-ui/data/licenseGroup.json
mongoimport /var/dfa/colreg-ui/data/accrualMethod.json
mongoimport /var/dfa/colreg-ui/data/accrualPolicy.json
mongoimport /var/dfa/colreg-ui/data/agentType.json
mongoimport /var/dfa/colreg-ui/data/language.json
# return success, even if service not installed
#exit 0
